# files/system/etc/profile.d/ublue-fastfetch.sh
alias neofetch='/usr/bin/fastfetch'
alias neowofetch='/usr/bin/fastfetch'
alias fastfetch='/usr/bin/fastfetch'
